# oldmess/services/server29/logwatch.service
[Unit]
Description=logwatcher service

[Service]
Type=simple
ExecStart=/root/bin/logwatcher

[Install]
WantedBy=multi-user.target
