# oldmess/services/server29/mailserver.service
[Unit]
Description=mailserver service

[Service]
Type=simple
ExecStart=/root/bin/sshtunnel mailserver.conf

[Install]
WantedBy=multi-user.target
