# services/mailserver.service
[Unit]
Description=mailserver service

[Service]
Type=simple
ExecStart=/root/bin/mailserver

[Install]
WantedBy=multi-user.target
